# deploy/systemd/room-monitor-firebase-collector.service
[Unit]
Description=Room Monitor Raspberry Pi Firebase collector
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/stm32_raspberrypi_practice/raspberry_pi
EnvironmentFile=/etc/room-monitor/firebase_collector.env
ExecStart=/opt/stm32_raspberrypi_practice/raspberry_pi/.venv/bin/python /opt/stm32_raspberrypi_practice/raspberry_pi/firebase_uploader.py --port ${SERIAL_PORT} --baud ${SERIAL_BAUD} --db ${LOCAL_DB} --retention-days ${LOCAL_RETENTION_DAYS} --prune-every-seconds ${LOCAL_PRUNE_EVERY_SECONDS} --firebase-database-url ${FIREBASE_DATABASE_URL} --upload-every-seconds ${FIREBASE_UPLOAD_EVERY_SECONDS}
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/stm32_raspberrypi_practice/raspberry_pi/data

[Install]
WantedBy=multi-user.target
